# root/usr/lib/systemd/system/frps.service
[Unit]
Description=Frp Server Service
After=network.target

[Service]
Type=simple
User=nobody
Restart=on-failure
RestartSec=5s
ExecStart=/program/frp/frps -c /program/frp/frps.ini

[Install]
WantedBy=multi-user.target
